# energy-tracker.service
[Unit]
Description=Koffein-Tracker API Server
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/koffein-tracker
EnvironmentFile=/opt/koffein-tracker/.env.local
ExecStart=/usr/bin/npm run server
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
